# crontab
# m h  dom mon dow   command
*/2 * * * * cd /home/blalor/xbee_dispatcher && ./upload_data.py
